# jenkins/etc/supervisord.d/jenkins.ini
[program:jenkins]
command=/usr/bin/java -jar /usr/lib/jenkins/jenkins.war
user=jenkins

